# dots/server/unraid/boot/config/plugins/ca.update.applications/plugin_update.cron
# Generated cron settings for plugin autoupdates
0 0 * * * /usr/local/emhttp/plugins/ca.update.applications/scripts/updateApplications.php >/dev/null 2>&1
