# dots/linux/hyprland/home/.config/systemd/user/restart-hypr-waybar.service
[Unit] 
Description="Restart hypr-waybar service"

[Service]
ExecStart=systemctl --user restart hypr-waybar.service
