# dots/linux/hyprland/home/.config/systemd/user/waybar-config.path
[Unit]
Description="Monitor the /etc/passwd file for changes"

[Path]
PathModified=/usr/local/share/waybar
Unit=restart-hypr-waybar.service

[Install]
WantedBy=hypr-waybar.service
