# dots/server/unraid/boot/config/plugins/dynamix/docker-update.cron
# Generated docker monitoring schedule:
10 0,12 * * * /usr/local/emhttp/plugins/dynamix.docker.manager/scripts/dockerupdate check &> /dev/null

