# dots/server/unraid/boot/config/plugins/dynamix/monitor.cron
# Generated system monitoring schedule:
*/1 * * * * /usr/local/emhttp/plugins/dynamix/scripts/monitor &> /dev/null

