# dots/server/unraid/boot/config/plugins/dynamix/mover.cron
# Generated mover schedule:
0 */1 * * * /usr/local/sbin/mover &> /dev/null

