# dots/server/unraid/boot/config/plugins/dynamix/plugin-check.cron
# Generated plugins version check schedule:
10 0,12 * * * /usr/local/emhttp/plugins/dynamix.plugin.manager/scripts/plugincheck &> /dev/null

