# dots/server/unraid/boot/config/plugins/dynamix/status-check.cron
# Generated array status check schedule:
20 */6 * * * /usr/local/emhttp/plugins/dynamix/scripts/statuscheck &> /dev/null

