# dots/server/unraid/boot/config/plugins/dynamix/unraid-check.cron
# Generated Unraid OS update check schedule:
11 0,12 * * * /usr/local/emhttp/plugins/dynamix.plugin.manager/scripts/unraidcheck &> /dev/null

