# dots/linux/hyprland/home/.config/aconfmgr/files/etc/ntpd.service.d/ntpd-ipv4.conf
[Service]
ExecStart=
ExecStart=/usr/bin/ntpd -4 -g -u ntp:ntp
