# dots/server/unraid/boot/config/plugins/community.applications/notification_scan.cron

# CRON for CA background scanning of applications
11 * * * * php /usr/local/emhttp/plugins/community.applications/scripts/notices.php > /dev/null 2>&1

